# supervisor/supervisord.ini
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock;

[supervisord]
logfile = supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = False
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp

[program:redis]
command=redis-server
autostart=true
autorestart=true

[program:celeryd]
command=celery worker -A sketchy.celery
directory=/Users/sbehrens/oss/Sketchy
user=sbehrens
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600

[program:gunicorn]
command=gunicorn sketchy:app -b 0.0.0.0:8000
directory=/Users/sbehrens/oss/Sketchy
user=sbehrens
autostart=true
autorestart=true
redirect_stderr=true
